# scripts/linux-services/rtk-mercury-gateway.service
[Unit]
Description=Mercury gateway (Mercury-fork agent + Discord bot)
After=network-online.target rtk-cloudflared.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mercury
ExecStart=%h/mercury/.venv/bin/mercury gateway run -v
Restart=always
RestartSec=10s
StartLimitIntervalSec=3600
StartLimitBurst=10
StandardOutput=append:%h/.mercury/logs/gateway.out.log
StandardError=append:%h/.mercury/logs/gateway.err.log

[Install]
WantedBy=default.target
